# crontab

*/5 * * * * root python /app/src/app.py >> /var/log/cron.log
